# gpio-motor-control-watcher.path
[Path]
PathModified=/home/user/transfer-table-2023

[Install]
WantedBy=multi-user.target